# ops/systemd/aps-scrape.service
[Unit]
Description=APS AI Tracker: scrape latest statements
Documentation=file:///home/ben/projects/aps-ai-tracker/cron-scrape.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ben/projects/aps-ai-tracker
Environment=PATH=/home/ben/.local/share/mise/shims:/home/ben/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ben/projects/aps-ai-tracker/cron-scrape.sh
TimeoutStartSec=2h
